# scripts/gb-chroot
#!/sbin/openrc-run

start() {
    for i in $BOOT_CHROOT
    do
        ebegin "Mounting chroot $i"
        gb-mount m "$i"
        eend $?
    done
}

stop() {
    find "$CHROOT_TOP" -mindepth 1 -maxdepth 1 -type d ! -iname 'tmpfs' -printf %P'\n' > "$RUN_DIR/tmp"
    while IFS= read -r i
    do
        if [ ! -e "$RUN_DIR"/"$i" ]
        then
            einfo "Skipping $i"
            continue
        fi

        COUNT=$(<"$RUN_DIR"/"$i")

        if [ -z "$COUNT" ]
        then
            COUNT=0
        fi

        if [ $COUNT -eq 0 ]
        then
            einfo "Skipping $i"
            continue
        fi

        if [ $COUNT -gt 1 ]
        then
            ewarn "Is $i still used? Unmounting anyway."
            echo "1" > "$RUN_DIR"/"$i"
        fi
        ebegin "Unmounting chroot $i"
        gb-mount u "$i"
        eend $?
    done < "$RUN_DIR/tmp"
    rm "$RUN_DIR/tmp"

    if [ -n "$USE_TMPFS" ];
    then
        ebegin "Unmounting tmpfs"
            umount -f "$CHROOT_TOP/tmpfs"
        eend $?
    fi
}
